# deploy/mcp_server.service
[Unit]
Description=Crypto Trading MCP Server (SSE)
After=network.target

[Service]
Type=simple
User=crypto_trader
WorkingDirectory=/opt/app
Environment="PYTHONUNBUFFERED=1"
Environment="USE_SECRET_MANAGER=true"
Environment="PROJECT_ID="
Environment="VERTEX_REGION=global"
ExecStart=/opt/app/venv/bin/python -m mcp_server.server
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
