# deploy/scheduler.service
[Unit]
Description=Crypto Trading System Scheduler
After=network.target

[Service]
Type=simple
User=crypto_trader
WorkingDirectory=/opt/app
Environment="PYTHONUNBUFFERED=1"
Environment="USE_SECRET_MANAGER=true"
Environment="PROJECT_ID="
Environment="VERTEX_REGION=global"
Environment="ENABLE_GCS_ARCHIVE=true"
ExecStart=/opt/app/venv/bin/python scheduler.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
